# deploy/systemd/kanban-worker.service
[Unit]
Description=Kanban autonomous worker (Claude Code)
After=kanban-api.service

[Service]
Type=oneshot
ExecStart=/root/.claude/hooks/kanban-worker.sh
