# deploy/systemd/kanban-worker.timer
[Unit]
Description=Kanban worker — poll backlog every 5 min

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=10s

[Install]
WantedBy=timers.target
